# probe/systemd/soil-probe.service
[Unit]
Description=Soil moisture probe (one reading)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=raspberry
Group=raspberry
SupplementaryGroups=gpio i2c
WorkingDirectory=/home/raspberry/Soil-Moisture-Probe/probe
EnvironmentFile=-/etc/soil-probe/env
ExecStart=/home/raspberry/Soil-Moisture-Probe/probe/venv/bin/python3 probes.py
