# probe/systemd/soil-probe.timer
[Unit]
Description=Soil moisture probe every 10 minutes

[Timer]
OnCalendar=*:0/10
AccuracySec=1min

[Install]
WantedBy=timers.target
